# deployment/htmltopdf.service
[Unit]
Description=HTML-to-PDF Converter API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/html-to-pdf
Environment="PATH=/opt/html-to-pdf/venv/bin"
ExecStart=/opt/html-to-pdf/venv/bin/gunicorn \
    --workers 4 \
    --bind 0.0.0.0:5000 \
    --timeout 120 \
    --access-logfile /var/log/html-to-pdf-access.log \
    --error-logfile /var/log/html-to-pdf-error.log \
    app:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target